# rpi/systemd/sqa-bme680.service.template
[Unit]
Description=Sleep Quality Sensors - BME680 Publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sqs-bme680.env
WorkingDirectory=/opt/sleep-quality-advisor/rpi
Environment=PYTHONPATH=/opt/sleep-quality-advisor/rpi/src
ExecStart=/usr/bin/env uv run python -m main
Restart=always
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target


